# io.service
[Unit]
Description=IO AI Orchestrator Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/io/packages/daemon/dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=IO_LOG_LEVEL=info
WorkingDirectory=/opt/io

[Install]
WantedBy=multi-user.target
